# gpg/config/systemd/user/gpg-agent-ssh.socket.d/environment.conf
# SPDX-FileCopyrightText: 2022 Emil Lundmark <emil@lndmrk.se>
# SPDX-License-Identifier: GPL-3.0-or-later

[Socket]
ExecStartPost=/usr/bin/systemctl --user set-environment \
    SSH_AUTH_SOCK=%t/gnupg/S.gpg-agent.ssh
