# environment/config/profile.d/10-less.sh
# SPDX-FileCopyrightText: 2019 Emil Lundmark <emil@lndmrk.se>
# SPDX-License-Identifier: GPL-3.0-or-later

export PAGER="less"
export LESS="--quit-if-one-screen \
             --LONG-PROMPT \
             --RAW-CONTROL-CHARS \
             --chop-long-lines"
