# roles/environment/files/profile.d/10-dircolors.sh
# SPDX-FileCopyrightText: 2021 Emil Lundmark <emil@lndmrk.se>
# SPDX-License-Identifier: GPL-3.0-or-later

if [ -z "${LS_COLORS}" ]; then
  eval "$(dircolors --sh "${HOME}/.dir_colors")"
fi
