# roles/environment/files/profile.d/10-python.sh
# SPDX-FileCopyrightText: 2021 Emil Lundmark <emil@lndmrk.se>
# SPDX-License-Identifier: GPL-3.0-or-later

export MYPY_CACHE_DIR="${XDG_CACHE_HOME:-${HOME}/.cache}/mypy"
export PYTHONPYCACHEPREFIX="${XDG_CACHE_HOME:-${HOME}/.cache}/python"
